# testdata/cron.d/valid
# vim: ft=crontab

SHELL = sh

TZ = Asia/Tokyo
@daily  echo hello

TZ = UTC
0 0 * * *  echo world

# duplicated task should be ignored
0 0 * * *    echo  world

ENABLE_USER_COLUMN = yes
@reboot   admin   cat%concron%initialized!
