# testdata/cron.d/valid
# vim: ft=crontab

SHELL = sh

TZ = Asia/Tokyo
@daily  root  echo hello

TZ = UTC
0 0 * * *  ec2-user  echo world

# duplicated task should be ignored
0 0 * * *    ec2-user    echo  world

TZ = "w h e r e ?"
@reboot   admin   cat%concron%initialized!
